# s6/server/run
#!/usr/bin/bash

( cd /spellbreak-server/spellbreak-server/ && python3 -m server.server )
